# tools/deployers/platform/roles/master/kubernetes/templates/kubernetes-server.target.j2
[Unit]
Description=Kubernetes Server target
Requires=kubernetes-scheduler.target kubernetes-controller.target kubernetes-apiserver.target

[Install]
WantedBy=multi-user.target
