# platform-deploy/roles/master/kubernetes/templates/kubernetes-apiserver.target.j2
[Unit]
Description=Kubernetes Apiserver target
Requires=kubernetes-apiserver.service kubernetes-apiserver-check.service

[Install]
WantedBy=multi-user.target
