# platform-deploy/roles/master/kubernetes/templates/kubernetes-scheduler.target.j2
[Unit]
Description=Kubernetes Scheduler target
Requires=kubernetes-scheduler.service kubernetes-scheduler-check.service

[Install]
WantedBy=multi-user.target
